# deployment/raspi3/zischr.target
[Unit]
Description=zischr meta service
After=graphical.target

[Install]
WantedBy=multi-user.target
