# deployment/raspi3/zischr.target
[Unit]
Description=zischr meta service
Requires=graphical.target network.target multi-user.target

